# rust-daily-backend.service
[Unit]
Description=rust-daily-backend
After=network-online.target
Wants=network-online.target

[Service]
User=www-data12
Group=www-data

WorkingDirectory=/var/www12/html
Environment=RUST_DAILY_ENV=prod
Environment=HOME=/var/www12
Environment=XDG_CONFIG_HOME=/var/www12/.config
Environment=XDG_DATA_HOME=/var/www12/.local/share
Environment=PATH=/var/www12/bin:/var/www12/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/var/www12/html/rust-daily-backend

Restart=on-failure
RestartSec=2

# Supervisor stopasgroup/killasgroup → control-group
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=15

# Logging (journald by default)
# StandardOutput=append:/var/log/rust-daily-backend/rust-daily-backend.log
# StandardError=inherit

# Security hardening
NoNewPrivileges=true
#PrivateTmp=true
PrivateDevices=true
#ProtectSystem=strict
#ProtectHome=true
ReadWritePaths=/var/www12
AmbientCapabilities=
CapabilityBoundingSet=
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
