# rust-daily-backend.service
[Unit]
Description=rust-daily-backend
After=network-online.target
Wants=network-online.target

[Service]
User=www-data12
Group=www-data

WorkingDirectory=/var/www12/html
Environment=RUST_DAILY_ENV=prod
Environment=HOME=/var/www12
Environment=XDG_CONFIG_HOME=/var/www12/.config
Environment=XDG_DATA_HOME=/var/www12/.local/share
Environment=XDG_RUNTIME_DIR=/run/rust-daily-backend
Environment=PATH=/var/www12/bin:/var/www12/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/var/www12/html/rust-daily-backend
TemporaryFileSystem=/var/www12/rust-daily-runs:size=2G,nr_inodes=200000,mode=1777
RuntimeDirectory=rust-daily-backend
RuntimeDirectoryMode=0700

Restart=on-failure
RestartSec=2

# Supervisor stopasgroup/killasgroup → control-group
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=15

# Logging (journald by default)
# StandardOutput=append:/var/log/rust-daily-backend/rust-daily-backend.log
# StandardError=inherit

# Security hardening.
# Rootless Podman needs setuid newuidmap/newgidmap; keep no-new-privileges
# inside the runner container instead of on the service process.
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/rust-daily-backend /var/www12/.cache /var/www12/.config /var/www12/.local/share
AmbientCapabilities=
LimitNOFILE=65535
LimitCORE=0
MemoryMax=1G
MemorySwapMax=0
CPUQuota=200%
TasksMax=512

[Install]
WantedBy=multi-user.target
